# Library/LaunchAgents/com.vishal.toggle_cli_apps_dark_mode.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.vishal.toggle_cli_apps_dark_mode</string>

        <key>KeepAlive</key>
        <true />

        <key>StandardErrorPath</key>
        <string>/Users/vkushwaha/.toggle_cli_apps_dark_mode_stderr.log</string>

        <key>StandardOutPath</key>
        <string>/Users/vkushwaha/.toggle_cli_apps_dark_mode_stdout.log</string>

        <key>ProgramArguments</key>
        <array>
            <string>/opt/homebrew/bin/dark-notify</string>
            <string>-c</string>
            <string>/Users/vkushwaha/Documents/vishal/projects/scripts/toggle_dark_mode.sh</string>
        </array>
    </dict>
</plist>
